# systemd/xhttp-enforcer.service
[Unit]
Description=XHTTP Manager Limits Enforcer (oneshot)
After=xray.service

[Service]
Type=oneshot
WorkingDirectory=/opt/xhttp-manager/addon
Environment=PYTHONPATH=/opt/xhttp-manager/addon
ExecStart=/opt/xhttp-manager/venv/bin/python /opt/xhttp-manager/addon/core/limits_enforcer.py
User=root
StandardOutput=journal
StandardError=journal
